# deploy/systemd/hyprl_track_record_daily.timer.example
[Unit]
Description=HyprL Track Record Daily Snapshot Timer

[Timer]
OnCalendar=*-*-* 20:10:00
Persistent=true

[Install]
WantedBy=timers.target
