# deploy/systemd/hyprl-crypto.timer
[Unit]
Description=Run HyprL Crypto signals hourly (24/7)

[Timer]
# Run hourly, 24/7 (crypto never sleeps)
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
